# .config/systemd/user/pacman-update-notify.service
[Unit]
Description=Notify about available pacman updates

[Service]
Type=oneshot
ExecStart=%h/.config/driftwm/scripts/pacman_update_notify.sh
